# WiFiChecker.cron
*/1 * * * * /home/pi/AmbientSense/WiFiChecker.sh > /home/pi/AmbientSense/WiFiChecker.log
